# supervisord.conf
[supervisord]
nodaemon=true ; Запускать supervisord в foreground, обязательно для Docker
logfile=/dev/null ; Не пишем общий лог супервизора (логи процессов будут отдельно)
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:tailscaled]
command=/usr/sbin/tailscaled --state=mem: --tun=userspace-networking ; Команда запуска tailscaled
autostart=true          ; Запускать автоматически при старте supervisord
autorestart=true        ; Перезапускать, если упадет
priority=10             ; Запускать одним из первых
stdout_logfile=/dev/stdout ; Направлять stdout процесса в stdout контейнера
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr ; Направлять stderr процесса в stderr контейнера
stderr_logfile_maxbytes=0
startsecs=5             ; Дать 5 секунд на запуск перед тем, как считать его успешным

[program:darkhttpd]
command=/usr/sbin/darkhttpd /www --port 8080 ; Команда запуска darkhttpd
autostart=true
autorestart=true
priority=100            ; Запускать после основных служб
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:microsocks]
command=/usr/local/bin/microsocks -p 1080 ; Команда запуска microsocks (БЕЗ AUTH)
; Если нужна аутентификация, используйте переменные окружения:
; command=bash -c "/usr/local/bin/microsocks -p 1080 -u \"${PROXY_USER}\" -P \"${PROXY_PASSWORD}\""
; environment=PROXY_USER="%(ENV_PROXY_USER)s",PROXY_PASSWORD="%(ENV_PROXY_PASSWORD)s"
autostart=true
autorestart=true
priority=90             ; Запускать после tailscaled, но до httpd
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0